# a-systemd.socket
[Socket]
ListenStream=/run/www-node.sk

[Install]
WantedBy=sockets.target

