# setup/etc-systemd-system_nodemanager.socket
[Socket]
ListenStream=/run/solobit.sk

[Install]
WantedBy=sockets.target

